# big-digits/show-energy.conf
[program:show-energy]
directory=/home/machinekit/atbristol-megadrawbz/big-digits
autostart=true
command=python /home/machinekit/atbristol-megadrawbz/big-digits/show_energy.py
user=root

